# systemd/mpf.service
[Unit]
Description=MPF

[Service]
User=mpfuser
WorkingDirectory=/home/mpfuser/mpf-data
ExecStart=mpf

[Install]
WantedBy=multi-user.target